# UDPServer/debian/postinst
#!/bin/bash

update-rc.d hddmonitor-udpserver defaults
