# LinuxSoftware/debian/postinst
#!/bin/bash

update-rc.d hddmonitor-udpclient defaults
